# phpnginx/build/php/wechatTempMsg.conf
[program:wechatTempMsg]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/shaka/artisan queue:work --queue=wechatTempMsg --tries=3 --daemon
user=root
autostart=true
autorestart=true
numprocs=4
stdout_logfile = /var/log/supervisor/supervisor_stdout.log
stderr_logfile = /var/log/supervisor/supervisor_stderr.log
